# systemd/aria-check-drift.timer
[Unit]
Description=ARIA — daily drift check at 02:00

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
