# systemd/aria-intraday.timer
[Unit]
Description=ARIA — intraday snapshot every 4h

[Timer]
OnCalendar=*-*-* 00,04,08,12,16,20:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
