# systemd/aria-suggest-automations.service
[Unit]
Description=ARIA — weekly automation suggestions (Ollama)
After=network-online.target aria-correlations.service

[Service]
Type=oneshot
WorkingDirectory=/home/justin/Documents/projects/ha-intelligence-hub
ExecStart=/bin/bash -c '. ~/.env && exec /home/justin/Documents/projects/ha-intelligence-hub/.venv/bin/aria suggest-automations'
StandardOutput=append:/home/justin/.local/log/aria.log
StandardError=append:/home/justin/.local/log/aria.log
MemoryHigh=3G
MemoryMax=4G
