# systemd/aria-hub.service
[Unit]
Description=ARIA — Adaptive Residence Intelligence Architecture
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=/home/justin/Documents/projects/ha-aria
ExecStart=/bin/bash -c '. ~/.env && exec /home/justin/Documents/projects/ha-aria/.venv/bin/aria serve --port 8001'
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=2G
CPUQuota=50%

[Install]
WantedBy=default.target
