# systemd/aria-sequences.service
[Unit]
Description=ARIA — weekly sequence training + anomaly detection
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/justin/Documents/projects/ha-intelligence-hub
ExecStart=/bin/bash -c '. ~/.env && /home/justin/Documents/projects/ha-intelligence-hub/.venv/bin/aria sequences train && /home/justin/Documents/projects/ha-intelligence-hub/.venv/bin/aria sequences detect'
StandardOutput=append:/home/justin/.local/log/aria.log
StandardError=append:/home/justin/.local/log/aria.log
MemoryHigh=3G
MemoryMax=4G
