# systemd/aria-prophet.service
[Unit]
Description=ARIA — weekly Prophet model training
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/justin/Documents/projects/ha-aria
ExecStart=/bin/bash -c '. ~/.env && exec /home/justin/Documents/projects/ha-aria/.venv/bin/aria prophet'
StandardOutput=append:/home/justin/.local/log/aria.log
StandardError=append:/home/justin/.local/log/aria.log
MemoryHigh=3G
MemoryMax=4G
