# systemd/aria-intraday.service
[Unit]
Description=ARIA — intraday snapshot
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/justin/Documents/projects/ha-intelligence-hub
ExecStart=/bin/bash -c '. ~/.env && exec /home/justin/Documents/projects/ha-intelligence-hub/.venv/bin/aria snapshot-intraday'
StandardOutput=append:/home/justin/.local/log/aria.log
StandardError=append:/home/justin/.local/log/aria.log
MemoryHigh=1G
MemoryMax=2G
